# dist-dry-run/work/cli-darwin-arm64/aegis-v1.1.0-darwin-arm64/packaging/systemd/aegis-edge-sitl.service
# Aegis Edge SITL example only. Do not install or enable automatically.
# Local SITL evidence only; no real-flight readiness claim.
[Unit]
Description=Aegis Edge SITL example

[Service]
Type=oneshot
ExecStart=/usr/local/bin/aegis-edge deployment check --profile /opt/aegis-edge/examples/edge/deployment/profiles/px4-sitl-local.yaml
