# dist-dry-run/work/cli-linux-arm64/aegis-v1.1.0-linux-arm64/packaging/systemd/aegis-edge-bench.example.service
# Aegis Edge bench example only. Do not install or enable automatically.
# hardware_bench_no_actuation only; no real-flight readiness claim.
[Unit]
Description=Aegis Edge bench example

[Service]
Type=oneshot
ExecStart=/usr/local/bin/aegis-edge bench doctor
