# packaging/systemd/edge.example.service
# Edge example only. Do not install, start, or enable automatically.
# Verify the release artifact checksum before copying binaries onto a host.
# Simulation/SITL/bench-preparation only; no real-flight or hardware endpoint.
[Unit]
Description=Edge example deployment diagnostics

[Service]
Type=oneshot
ExecStart=/usr/local/bin/edge deployment doctor
