# packaging/systemd/edge-sitl.service
# Edge SITL example only. Do not install or enable automatically.
# Local SITL evidence only; no real-flight readiness claim.
[Unit]
Description=Edge SITL example

[Service]
Type=oneshot
ExecStart=/usr/local/bin/edge deployment check --profile /opt/edge/examples/edge/deployment/profiles/px4-sitl-local.yaml
